# files/xapsd.in
#!/bin/sh
# $FreeBSD$
#

# PROVIDE: xapsd
# REQUIRE: LOGIN cleanvar dovecot
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf[.local] to enable xaps:
# xapsd_enable="YES"
#

. /etc/rc.subr

name="xapsd"
rcvar=${name}_enable

load_rc_config $name

#defaults
: ${xapsd_enable="NO"}
: ${xapsd_cmd="/usr/local/bin/xapsd"}
: ${xapsd_pidfile="/var/run/xapsd/xapsd.pid"}
: ${xapsd_log="/var/log/xapsd.log"}

pidfile=${xapsd_pidfile}
command="/usr/sbin/daemon"
command_args="-P ${pidfile} -r -f -m 3 -o ${xapsd_log} ${xapsd_cmd} ${xapsd_args}"

run_rc_command "$1"
